# ops/com.zhuanz1.cloud-hub-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zhuanz1.cloud-hub-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/Zhuanz1/Documents/solo-template-business &amp;&amp; /usr/bin/python3 /Users/Zhuanz1/Documents/solo-template-business/scripts/sync_cloud_hub_to_local.py --pull --trigger-label launchd</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/Zhuanz1/Documents/solo-template-business</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/Zhuanz1/Documents/solo-template-business/data/automation_logs/cloud_hub_sync.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/Zhuanz1/Documents/solo-template-business/data/automation_logs/cloud_hub_sync.stderr.log</string>
</dict>
</plist>
